# services/tns/run
#!/bin/sh
exec java -jar /tugboat-event-receiver.jar
